# packages/agent/platform/darwin/installer/LaunchDaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- SMAppService daemon plist. Ships INSIDE the app bundle at
         Contents/Library/LaunchDaemons/com.nexus-gateway.agent.plist and is
         registered into the system domain by the menu app via
         SMAppService.daemon(plistName:). It is NOT installed to
         /Library/LaunchDaemons — the registration is bundle-tied, so deleting
         the app deregisters the daemon (no orphaned plist, no launchctl
         bootstrap, no disabled-override EIO). See LaunchServiceManager.swift. -->

    <!-- Label must match the plist file name's reverse-DNS stem. -->
    <key>Label</key>
    <string>com.nexus-gateway.agent</string>

    <!-- BundleProgram: the daemon executable relative to the app bundle root.
         Replaces the absolute /Applications/... path the classic plist used —
         SMAppService resolves it relative to whichever bundle registered it, so
         the daemon follows the app even if installed to a non-standard path. -->
    <key>BundleProgram</key>
    <string>Contents/MacOS/nexus-agent</string>

    <!-- argv for the daemon. argv[0] is conventional; the run subcommand + the
         config path (still installed to /Library/Application Support by the
         pkg postinstall) are the load-bearing arguments. -->
    <key>ProgramArguments</key>
    <array>
        <string>nexus-agent</string>
        <string>run</string>
        <string>-config</string>
        <string>/Library/Application Support/com.nexus-gateway.agent/agent.yaml</string>
    </array>

    <!-- Boot daemon: start at load so enforcement is active before any GUI
         login. SMAppService.daemon registers a system-domain LaunchDaemon
         (root), so this preserves pre-login enforcement — we do NOT downgrade
         to a per-user login agent. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Run as root. SMAppService.daemon registers a system-domain daemon that
         runs as root by default, but the daemon hard-requires it (the NE status
         socket path and all /Library state are root-owned, and the peer-cred
         check expects root), so we state it explicitly rather than rely on the
         default. -->
    <key>UserName</key>
    <string>root</string>

    <!-- KeepAlive is load-bearing for the user-quit-flag coupling: the menu's
         Quit writes a flag file the daemon reads on every respawn and exits
         while present; KeepAlive respawns it (throttled) and re-opening the app
         clears the flag. Each respawn also re-reads the quitAllowed policy.
         Must survive the SMAppService plist key subset (validate on-device). -->
    <key>KeepAlive</key>
    <true/>

    <!-- Capture stdout/stderr for debugging. SMAppService may constrain where a
         daemon can write; if these /Library/Logs paths are rejected, re-home
         the logs (validate on-device — see the P2 design plist-subset task). -->
    <key>StandardOutPath</key>
    <string>/Library/Logs/com.nexus-gateway.agent/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/com.nexus-gateway.agent/launchd.err.log</string>

    <!-- Resource limits. -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <!-- Throttle the respawn loop (also bounds the quit-flag respawn-and-die
         cadence). Load-bearing for P1 — must survive the key subset. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
