# packages/agent/platform/darwin/installer/LaunchAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Per-user agent: launchd runs this for every user that logs in.
         Sole job is to (re-)open NexusAgent.app at session start so the
         menu-bar app gets a chance to re-issue the NETransparentProxy
         startVPNTunnel call after a reboot. The system LaunchDaemon
         (com.nexus-gateway.agent) already starts the Go daemon at boot,
         but only a per-user GUI app can call the NetworkExtension
         framework APIs that actually wire macOS → our provider. Without
         this agent the user has to manually open NexusAgent.app after
         every reboot before traffic interception begins. -->
    <key>Label</key>
    <string>com.nexus-gateway.agent.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <!-- -g: do NOT bring the app to the foreground. NexusAgent.app is
             LSUIElement=true so it has no Dock icon either way; -g is belt
             and suspenders so the app starts silently in the menu bar
             without grabbing focus from whatever the user is doing. -->
        <string>-g</string>
        <string>/Applications/NexusAgent.app</string>
    </array>

    <!-- Run at every user login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- LimitLoadToSessionType: only run inside a normal GUI login session.
         Without this filter launchd would also try to launch the agent at
         pre-login (loginwindow) where /Applications/NexusAgent.app cannot
         render a menu bar item, producing a tight respawn loop. -->
    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>

    <!-- One-shot semantics: macOS calls `open` once, `open` exits as soon
         as the .app is launched (the .app itself is a long-running
         process). We do NOT want launchd to keep respawning `open` — it
         would chew CPU and `open` against an already-running .app is a
         no-op anyway. -->
    <key>KeepAlive</key>
    <false/>

    <!-- ThrottleInterval prevents respawn storms if the user closes
         NexusAgent.app immediately after login. With KeepAlive=false this
         is mostly defense-in-depth. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
